# supervisord.conf
[supervisord]
#This Should be change application owner
[program:material-api]
command=/bin/bash -c "node /usr/src/app/index.js"
startsecs = 0
autorestart = unexpected
exitcodes = 0