# supervisord.conf
[supervisord]
#This Should be change application owner
[program:material-api]
command=/bin/bash -c "node index.js"
startsecs = 0
autorestart = unexpected
exitcodes = 0